# cc-ms-artblocks-nft/deploy/systemd/cc-ms-artblocks-nft.timer
[Unit]
Description=Delay timer for cc-ms-artblocks-nft.service

[Timer]
OnBootSec=10

[Install]
WantedBy=timers.target